# systemd/root/mouseless.service
[Unit]
Description=mouseless

[Service]
ExecStartPre=/bin/sleep 2
ExecStart=/usr/local/bin/mouseless --config /home/anthony/.config/mouseless/config.yaml

[Install]
WantedBy=multi-user.target
